# Plex-Media-Server-SYNC/crontab/delete_local_cloud.sh
#!/bin/bash

rm -r "/media/Plex_local/Camera_Current/"
rm -r "/media/Plex_local/Camera_Phone/"
rm -r "/media/Plex_local/Music/
rm -r "/media/Plex_local/Pictures/"
rm -r "/media/Plex_local/Videos/"
rm -r "/media/Plex_local/Movies/"
rm -r "/media/Plex_local/TV_Shows/"
rm -r "/media/Plex_local/Reactions/"
